# startupplay.service
[Unit]
Description=Startup Play
After=multi-user.target

[Service]
Type=idle
ExecStart=/etc/startupplay.sh

[Install]
WantedBy=multi-user.target
